# root/etc/services.d/php7-fpm/run
#!/usr/bin/execlineb -P
php-fpm7 -c /config/php7/php.ini --fpm-config /config/php7/php-fpm.conf -F
