# root/etc/services.d/nginx/run
#!/usr/bin/execlineb -P
nginx -c /config/nginx/nginx.conf -g "daemon off;"
